# .config/autostart/check-update.desktop
[Desktop Entry]
Type=Application
Name=Check System Update
Exec=/usr/bin/bash -c "$HOME/.local/bin/check-update"
StartupNotify=false
