# .config/autostart/remove-screenshots.desktop
[Desktop Entry]
Type=Application
Name=Remove Screenshots
Exec=/usr/bin/bash -c "rm ~/Pictures/Screenshots/*"
Terminal=false
StartupNotify=false
