# .config/systemd/user/syncthing.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/bash -c 'until findmnt --source UUID=a85a18b8-85fd-4107-a2fe-cd18ba52f181 --mountpoint /mnt/storage0 >/dev/null; do sleep 1; done'
